# cron/crontab
*/2 * * * * cd /code && python scripts/backup_secure.py >> /tmp/secure_backup.log 2>&1
